# clash-server
#!/sbin/openrc-run

name="clash"
command="clash -d /root/.config/"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
   need net
}

start() {
   ebegin "Starting Clash"
   start-stop-daemon --start --background \
        --exec $command \
        --make-pidfile --pidfile $pidfile
   eend $?
}

stop() {
   ebegin "Stopping Clash"
   start-stop-daemon --stop \
        --exec $command \
        --pidfile $pidfile
   eend $?
}
